# conf/supervisor_nuvem.conf
[program:nuvem]
command = gunicorn nuvem:app
directory = /home/pi/nuvem/
user = pi

[program:homebridge]
command = homebridge
directory = /home/pi/
user = pi